# server/deploy_configs/supervisor/production.ini
[program:{{ project_name }}-production]
directory = <path to {{ project_name }} root directory>
command = <path to virtualenv uwsgi binary> --ini <path to {{project_name}} root directory>/deploy_configs/uwsgi/production.ini
autostart = true
autorestart = true

