# sc_server_django/server/supervisor/production.ini
[program:studycentric-production]
directory = /home/devuser/webapps/studycentric-env/studycentric/sc_server_django
command = /home/devuser/webapps/studycentric-env/bin/uwsgi --ini /home/devuser/webapps/studycentric-env/studycentric/sc_server_django/server/uwsgi/production.ini
autostart = true
autorestart = true
