# sc_server_django/server/supervisor/production.ini
[program:studycentric-production]
directory = /home/devuser/sites/studycentric-env/django-server
command = /home/devuser/sites/studycentric-env/bin/uwsgi --ini /home/devuser/sites/studycentric-env/django-server/server/uwsgi/production.ini
autostart = true
autorestart = true
